# scripts/deesweets-webapp.service
[Unit]
Description=Dee Sweets Analytics Web Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/DeeSweetsAnalytics
Environment="PATH=/home/ubuntu/DeeSweetsAnalytics/webapp/.venv/bin"
Environment="GOOGLE_APPLICATION_CREDENTIALS=/home/ubuntu/websitecountryspikes-c44a6b026c7b.json"
Environment="APP_HOST=0.0.0.0"
Environment="APP_PORT=5000"
ExecStart=/home/ubuntu/DeeSweetsAnalytics/webapp/.venv/bin/python webapp/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

